# etc/zuul-scheduler.service.d/centos.conf
[Service]
ExecStart=
ExecStart=/usr/bin/zuul-scheduler
ExecStop=
ExecStop=/usr/bin/zuul-scheduler stop
