# etc/zuul-web.service.d/centos.conf
[Service]
ExecStart=
ExecStart=/usr/bin/zuul-web
