# etc/zuul-executor.service.d/centos.conf
[Service]
ExecStart=
ExecStart=/usr/bin/zuul-executor
ExecStop=
ExecStop=/usr/bin/zuul-executor stop
